# dot_config/systemd/user/ollama.service
[Unit]
Description=Ollama API server
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/ollama serve
Restart=always
RestartSec=2
# Uncomment to customize:
# Environment=OLLAMA_HOST=127.0.0.1:11434
# Environment=OLLAMA_MODELS=%h/.ollama/models
# Environment=OLLAMA_NUM_PARALLEL=1
#
Environment=OLLAMA_USE_CUDA=1 
#Environment=OLLAMA_DEBUG=1
#

[Install]
WantedBy=default.target

